# etcd/etcdstart.sh
#!/bin/sh
#chkconfig: 2345 80 90
#description:etcd_开机自启脚本

#开启etcd服务
